# packages/csr_ha/csr_ha-2.0.2.tar.gz/csr_ha-2.0.2/csr_ha.service
[Unit]
Description=CSR High Availability service
After=waagent.service network.target

[Service]
Type=simple
User=guestshell
Group=guestshell
WorkingDirectory=/home/guestshell/cloud/HA
ExecStart=/home/guestshell/.local/lib/python2.7/site-packages/csr_ha/server/ha_server.py start
Restart=always
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
